# kasandra-crawler/scripts/kasCrawler.service
[Unit]
Description=Kasandra Crawler
[Service]
ExecStart=/data/kasandra/runCrawler.sh
Type=simple
User=root
[Install]
WantedBy=multi-user.target
